# tools/windy-render.service
[Unit]
Description=Windy display – render weather image
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Run as root or a dedicated user; adjust User= if preferred
WorkingDirectory=/opt/windy
ExecStart=/usr/bin/python3 /opt/windy/windy_render.py \
    --out-bin-temp /var/lib/windy/windy_temp.bin \
    --out-bin-hum  /var/lib/windy/windy_hum.bin  \
    --out-png-temp /var/lib/windy/windy_temp.png \
    --out-png-hum  /var/lib/windy/windy_hum.png  \
    --no-header
StandardOutput=journal
StandardError=journal
